# com.shadowrocket.autoconfig.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shadowrocket.autoconfig</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mac/My/2-Research/terminal-work/Shadowrocket-Russia-config/autoconfig/monitor.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/shadowrocket-autoconfig.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/shadowrocket-autoconfig.err</string>
</dict>
</plist>
